# pipeline/supervisord-luigi.ini
[program:luigid]
environment = PYTHONPATH=$PYTHONPATH:/opt/pipeline/tasks/
command=luigid --background --pidfile /var/run/luigid.pid --logdir /opt/pipeline/logs